# hay_cm5_fan/rootfs/etc/services.d/cm5-fan/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: CM5 Fan Controller
# Take down the S6 supervision tree when CM5 Fan Controller fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "CM5 Fan Controller crashed with exit code ${1}. Respawning..."
fi
