# sonuntius/rootfs/etc/services.d/tidal-connect/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: Sonuntius
# tidal-connect finish handler.
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "tidal-connect exited with code ${1}. Respawning..."
fi
