# sonuntius/rootfs/etc/services.d/yt-cast/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: Sonuntius
# Take down the S6 supervision tree when yt-cast fails.
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "Sonuntius yt-cast crashed with exit code ${1}. Respawning..."
fi
