# sonuntius/rootfs/etc/services.d/yt-cast/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: Sonuntius
# YouTube / YT Music DIAL + lounge receiver — Go port of yt-cast-receiver.
# Binary at /usr/local/bin/yt-cast.
# ==============================================================================
bashio::log.info "Starting Sonuntius yt-cast..."
exec /usr/local/bin/yt-cast
