# muninndb/rootfs/etc/services.d/muninndb-provision/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: MuninnDB
# Provisioning service finish handler — do not restart on exit
# ==============================================================================
