# muninndb/rootfs/etc/services.d/muninndb-provision/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: MuninnDB
# Post-startup provisioning: sets admin password and creates default vault.
# Runs as an S6 service that exits after provisioning is complete.
# ==============================================================================

PROVISION_MARKER="/data/muninndb/.provisioned"
STORED_HASH_FILE="/data/muninndb/.admin_pass_hash"
STORED_PASS_FILE="/data/muninndb/.admin_pass_set"

# If already provisioned and no password change requested, exit immediately
if [[ -f "${PROVISION_MARKER}" ]]; then
    if ! bashio::config.has_value 'admin_password' || [[ -z "$(bashio::config 'admin_password')" ]]; then
        bashio::log.info "MuninnDB already provisioned, skipping"
        exec s6-svc -O /run/service/muninndb-provision
    fi
fi

# Wait for MuninnDB to be healthy
bashio::log.info "Waiting for MuninnDB to start..."
MAX_WAIT=60
WAITED=0
while [[ ${WAITED} -lt ${MAX_WAIT} ]]; do
    if curl -sf http://127.0.0.1:8476/api/health > /dev/null 2>&1; then
        bashio::log.info "MuninnDB is ready"
        break
    fi
    sleep 1
    WAITED=$((WAITED + 1))
done

if [[ ${WAITED} -ge ${MAX_WAIT} ]]; then
    bashio::log.error "MuninnDB did not become healthy within ${MAX_WAIT}s, skipping provisioning"
    exec s6-svc -O /run/service/muninndb-provision
fi

# Determine current password (use default on first run, or stored password)
CURRENT_PASSWORD="password"
if [[ -f "${STORED_PASS_FILE}" ]]; then
    CURRENT_PASSWORD=$(cat "${STORED_PASS_FILE}")
fi

# Login to get session cookie
SESSION_COOKIE=$(mktemp)
LOGIN_PAYLOAD=$(jq -n --arg user "root" --arg pass "${CURRENT_PASSWORD}" \
    '{username: $user, password: $pass}')

LOGIN_RESPONSE=$(curl -s -w "\n%{http_code}" -c "${SESSION_COOKIE}" \
    -X POST http://127.0.0.1:8476/api/auth/login \
    -H 'Content-Type: application/json' \
    -d "${LOGIN_PAYLOAD}" 2>/dev/null)

HTTP_CODE=$(echo "${LOGIN_RESPONSE}" | tail -1)
if [[ "${HTTP_CODE}" != "200" ]]; then
    LOGIN_BODY=$(echo "${LOGIN_RESPONSE}" | sed '$d')
    bashio::log.warning "Could not login to MuninnDB (HTTP ${HTTP_CODE}): ${LOGIN_BODY}"
    rm -f "${SESSION_COOKIE}"
    exec s6-svc -O /run/service/muninndb-provision
fi

bashio::log.info "Logged in to MuninnDB as root"

# Change admin password if configured
if bashio::config.has_value 'admin_password' && [[ -n "$(bashio::config 'admin_password')" ]]; then
    NEW_PASSWORD=$(bashio::config 'admin_password')
    NEW_HASH=$(printf '%s' "${NEW_PASSWORD}" | sha256sum | cut -d' ' -f1)

    # Compare against stored hash (not the password text)
    STORED_HASH=""
    if [[ -f "${STORED_HASH_FILE}" ]]; then
        STORED_HASH=$(cat "${STORED_HASH_FILE}")
    fi

    if [[ ${#NEW_PASSWORD} -lt 8 ]]; then
        bashio::log.error "admin_password must be at least 8 characters (got ${#NEW_PASSWORD})"
    elif [[ "${NEW_HASH}" == "${STORED_HASH}" ]]; then
        bashio::log.info "Admin password unchanged"
    else
        bashio::log.info "Password change detected, updating..."

        PASS_PAYLOAD=$(jq -n --arg user "root" --arg pass "${NEW_PASSWORD}" \
            '{username: $user, new_password: $pass}')

        PASS_RESPONSE=$(curl -s -w "\n%{http_code}" -b "${SESSION_COOKIE}" \
            -X PUT http://127.0.0.1:8476/api/admin/password \
            -H 'Content-Type: application/json' \
            -d "${PASS_PAYLOAD}" 2>/dev/null)

        PASS_CODE=$(echo "${PASS_RESPONSE}" | tail -1)
        PASS_BODY=$(echo "${PASS_RESPONSE}" | sed '$d')

        if [[ "${PASS_CODE}" == "200" ]] || [[ "${PASS_CODE}" == "204" ]]; then
            bashio::log.info "Admin password updated successfully"
            # Store both the hash (for comparison) and the password (for login)
            printf '%s' "${NEW_HASH}" > "${STORED_HASH_FILE}"
            chmod 600 "${STORED_HASH_FILE}"
            printf '%s' "${NEW_PASSWORD}" > "${STORED_PASS_FILE}"
            chmod 600 "${STORED_PASS_FILE}"

            # Re-login with new password for vault creation
            RELOGIN_PAYLOAD=$(jq -n --arg user "root" --arg pass "${NEW_PASSWORD}" \
                '{username: $user, password: $pass}')
            curl -s -c "${SESSION_COOKIE}" \
                -X POST http://127.0.0.1:8476/api/auth/login \
                -H 'Content-Type: application/json' \
                -d "${RELOGIN_PAYLOAD}" > /dev/null 2>&1
        else
            bashio::log.warning "Failed to update admin password (HTTP ${PASS_CODE}): ${PASS_BODY}"
        fi
    fi
fi

# Ensure default vault exists
if bashio::config.has_value 'default_vault' && [[ -n "$(bashio::config 'default_vault')" ]]; then
    VAULT_NAME=$(bashio::config 'default_vault')

    # GET /api/vaults returns a JSON array of vault name strings: ["default","homeassistant"]
    VAULT_LIST=$(curl -s -b "${SESSION_COOKIE}" \
        "http://127.0.0.1:8476/api/vaults" 2>/dev/null)

    if echo "${VAULT_LIST}" | jq -e --arg name "${VAULT_NAME}" \
        'index($name) != null' > /dev/null 2>&1; then
        bashio::log.info "Vault '${VAULT_NAME}' already exists"
    else
        VAULT_PUBLIC=true
        if ! bashio::config.true 'default_vault_public'; then
            VAULT_PUBLIC=false
        fi

        VAULT_PAYLOAD=$(jq -n --arg name "${VAULT_NAME}" --argjson public "${VAULT_PUBLIC}" \
            '{name: $name, public: $public}')

        VAULT_RESPONSE=$(curl -s -w "\n%{http_code}" -b "${SESSION_COOKIE}" \
            -X PUT http://127.0.0.1:8476/api/admin/vaults/config \
            -H 'Content-Type: application/json' \
            -d "${VAULT_PAYLOAD}" 2>/dev/null)

        VAULT_CODE=$(echo "${VAULT_RESPONSE}" | tail -1)
        if [[ "${VAULT_CODE}" == "200" ]] || [[ "${VAULT_CODE}" == "201" ]] || [[ "${VAULT_CODE}" == "204" ]]; then
            if [[ "${VAULT_PUBLIC}" == "true" ]]; then
                bashio::log.info "Created vault '${VAULT_NAME}' (public)"
            else
                bashio::log.info "Created vault '${VAULT_NAME}' (private — API key required)"
            fi
        else
            VAULT_BODY=$(echo "${VAULT_RESPONSE}" | sed '$d')
            bashio::log.warning "Failed to create vault '${VAULT_NAME}' (HTTP ${VAULT_CODE}): ${VAULT_BODY}"
        fi
    fi
fi

# Mark as provisioned
touch "${PROVISION_MARKER}"

# Clean up
rm -f "${SESSION_COOKIE}"

bashio::log.info "MuninnDB provisioning complete"

# Exit this oneshot service
exec s6-svc -O /run/service/muninndb-provision
