# dockhand/rootfs/etc/services.d/dockhand/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: Dockhand
# Take down the S6 supervision tree when Dockhand fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "Dockhand crashed with exit code ${1}. Respawning..."
fi
