# aegis_ha/rootfs/etc/services.d/aegis_ha/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: AegisHA
# S6 finish handler. A clean SIGTERM shutdown exits 0; the S6 supervisor
# reports signal exits as 256. Anything else is an unexpected crash.
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "AegisHA crashed with exit code ${1}. Respawning..."
fi
