# sonuntius/rootfs/etc/services.d/cast-receiver/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: Sonuntius
# Cast (CASTV2 + mDNS) receiver service.
# Binary at /usr/local/bin/cast-receiver. Reads /data/options.json
# directly; tolerates a missing AirReceiver cert by suppressing the TLS
# listener and staying alive so S6 does not enter a restart loop.
# ==============================================================================
bashio::log.info "Starting Sonuntius cast-receiver..."
exec /usr/local/bin/cast-receiver
