# huly/rootfs/etc/services.d/huly/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Huly
# Starts the Huly stack via Docker Compose
# ==============================================================================
bashio::log.info "========================================"
bashio::log.info "Huly stack starting at $(date '+%Y-%m-%d %H:%M:%S %Z')"
bashio::log.info "========================================"

cd /data/huly || exit 1

bashio::log.debug "Working directory: $(pwd)"
bashio::log.debug "Compose file exists: $(test -f compose.yml && echo 'yes' || echo 'no')"
bashio::log.debug "Env file exists: $(test -f .env && echo 'yes' || echo 'no')"

# Verify nginx.conf is a file (not a directory from a previous failed mount)
if [[ -d /data/huly/nginx.conf ]]; then
    bashio::log.error "nginx.conf is a directory instead of a file — cleaning up stale Docker mount artifact..."
    rm -rf /data/huly/nginx.conf
    # Re-run init to regenerate the config file
    bashio::log.info "Regenerating nginx configuration..."
    /etc/cont-init.d/huly.sh
fi
if [[ ! -f /data/huly/nginx.conf ]]; then
    bashio::log.error "nginx.conf does not exist — cannot start Huly stack"
    exit 1
fi

# Pull images if needed (first run or version change)
CURRENT_VERSION="${HULY_VERSION:-unknown}"
PULLED_VERSION=""
if [[ -f /data/huly/.images_pulled ]]; then
    PULLED_VERSION=$(cat /data/huly/.images_pulled 2>/dev/null)
fi

if [[ "${PULLED_VERSION}" != "${CURRENT_VERSION}" ]]; then
    if [[ -z "${PULLED_VERSION}" ]]; then
        bashio::log.info "Pulling Huly images (first run, this may take a while)..."
    else
        bashio::log.info "Huly version changed (${PULLED_VERSION} -> ${CURRENT_VERSION}), pulling new images..."
    fi
    bashio::log.debug "Running: docker-compose pull"
    if /usr/local/bin/docker-compose pull; then
        echo "${CURRENT_VERSION}" > /data/huly/.images_pulled
        bashio::log.info "Images pulled successfully"
    else
        bashio::log.error "Failed to pull Huly images"
        bashio::log.debug "docker-compose pull exit code: $?"
        exit 1
    fi
else
    bashio::log.debug "Images already pulled for version ${CURRENT_VERSION}"
fi

# Start the stack in detached mode so compose container logs don't flood HA.
# S6 needs a foreground process, so we keep alive with a health monitor loop.
bashio::log.info "Starting Huly services..."
/usr/local/bin/docker-compose up -d --remove-orphans 2>&1

if [[ $? -ne 0 ]]; then
    bashio::log.error "docker-compose up failed"
    exit 1
fi

bashio::log.info "Huly services started in background"
bashio::log.info "Use 'docker compose -f /data/huly/compose.yml logs -f' on the host to view service logs"

# Health monitor — keeps S6 service alive and periodically checks compose stack
HEALTH_INTERVAL=30
while true; do
    sleep "${HEALTH_INTERVAL}"

    # Count running compose containers
    RUNNING=$(/usr/local/bin/docker-compose ps --status running -q 2>/dev/null | wc -l)
    TOTAL=$(/usr/local/bin/docker-compose ps -q 2>/dev/null | wc -l)

    if [[ "${RUNNING}" -eq 0 ]]; then
        bashio::log.error "All Huly containers have stopped unexpectedly"
        exit 1
    fi

    if [[ "${RUNNING}" -lt "${TOTAL}" ]]; then
        # Log which containers are down
        DOWN=$(/usr/local/bin/docker-compose ps --status exited --format '{{.Name}}' 2>/dev/null | tr '\n' ', ')
        bashio::log.warning "Huly: ${RUNNING}/${TOTAL} containers running (down: ${DOWN%,})"
    fi
done
